# onboot_script.conf
start on (started networking)

script
  /usr/local/bin/onboot_script
end script
